# cleanup-tmp.service
[Unit]
Description=Service zum Löschen temporärer Dateien

[Service]
Type=oneshot
ExecStart=/bin/rm -rf /tmp/*

[Install]
WantedBy=multi-user.target
