# cleanup-tmp.timer
[Unit]
Description=Timer zum Löschen temporärer Dateien

[Timer]
OnCalendar=*-*-* 02:00:00
Persistent=true

[Install]
WantedBy=timers.target
